# config/etc/supervisor/conf.d/nodejs.conf
[program:nodejs]
directory=/ghost
command=/usr/bin/npm start
stopsignal=QUIT
redirect_stderr=true
autorestart=true
