# config/etc/supervisor/conf.d/wait-generate-thumbs.conf
[program:wait-generate-thumbs]
command=/usr/local/bin/wait-generate-thumbs
stopsignal=TERM
redirect_stderr=true
autorestart=true
